# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:proprio]
priority=10
directory=/opt/apps/proprio
command=/usr/local/bin/uwsgi
    --http-socket 0.0.0.0:8000
    -p 4
    -b 32768
    -T
    --master
    --max-requests 5000
    -H /opt/ve/proprio
    --static-map /static=/opt/apps/proprio/static
    --module wsgi:application
user=proprio
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
